# etc/skel/.config/autostart/compton.desktop
[Desktop Entry]
Exec=compton
Name=compton
Type=Application
Version=1.0
