# etc/skel/.config/autostart/ksuperkey.desktop
[Desktop Entry]
Exec=ksuperkey
Name=ksuperkey
Type=Application
Version=1.0
