# etc/skel/.config/autostart/picom.desktop
[Desktop Entry]
Exec=picom -b --config ~/.config/lxqt/picom.conf
Name=picom
Type=Application
Version=1.0
