# etc/skel/.config/autostart/plank.desktop
[Desktop Entry]
Exec=plank
Name=plank
Type=Application
Version=1.0
